# docusaurus.service
[Unit]
Description=Docusaurus Server for Physical AI Textbook
After=network.target

[Service]
Type=simple
User=seherkhan
WorkingDirectory=/mnt/c/Users/user/Desktop/my_ai_book
ExecStart=/usr/bin/npm run serve:prod
Restart=always
RestartSec=10
Environment=NODE_ENV=production
Environment=PORT=3000

[Install]
WantedBy=multi-user.target